# FILES.cluster-backup/etc/cron.d/backup
22 5 * * * root /usr/local/bin/backup >> /var/log/cron.log 2>&1


