# FILES.cluster-backup/etc/cron.d/test
* * * * * root echo "hello cron" >> /var/log/cron.log 2>&1


