# crontab
# Run backup script nightly at 5:00AM UTC
0 5 * * * root env - $(cat /root/env.sh) /root/backup.sh >> /var/log/cron.log 2>&1


