# crontab
22 5 * * * root /root/backup.sh >> /var/log/cron.log 2>&1


