# FILES.cluster-backup/var/spool/cron/crontabs/root
22 5 * * * /usr/local/bin/backup >> /var/log/cron.log 2>&1
